# nextjs-app/services/tts-watchdog.timer
[Unit]
Description=Run TTS watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
AccuracySec=30s

[Install]
WantedBy=timers.target
